# systemd/system/webkde.service.in
[Unit]
Description=WebKDE host-native Plasma streaming service
Documentation=file:/opt/webkde/README.md
Requires=docker.service user@@UID@.service
After=docker.service user@@UID@.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=WEBKDE_ENV_FILE=/etc/webkde/webkde.env
ExecStart=/opt/webkde/scripts/system-start.sh
ExecStop=/opt/webkde/scripts/system-stop.sh
TimeoutStartSec=20min
TimeoutStopSec=3min

[Install]
WantedBy=multi-user.target
