# systemd/user/webkde-session.service.in
[Unit]
Description=Host-native KDE Plasma session streamed by WebKDE

[Service]
Type=simple
Environment=WEBKDE_ENV_FILE=/etc/webkde/webkde.env
WorkingDirectory=/opt/webkde
ExecStartPre=/opt/webkde/scripts/wait-wayland.sh @RUNTIME@/wayland-0
ExecStart=/opt/webkde/scripts/start-plasma.sh
ExecStartPost=/opt/webkde/scripts/initial-display-mode.sh
ExecStartPost=@SYSTEMCTL@ --user start webkde-bridge.service
ExecStartPost=@SYSTEMCTL@ --user start webkde-inhibit.service
# startplasma exits cleanly when its outer Wayland connection disappears.
# Restart even after that clean exit so a transient container/Labwc restart
# cannot leave Selkies streaming an empty black compositor.
Restart=always
RestartSec=5s
TimeoutStartSec=5min
TimeoutStopSec=45s
