# systemd/user/webkde-bridge.service.in
[Unit]
Description=WebKDE clipboard and virtual-screen bridge
After=plasma-workspace.target
PartOf=webkde-session.service

[Service]
Type=simple
Environment=WEBKDE_ENV_FILE=@ENV_FILE@
ExecStart=@APP_DIR@/scripts/webkde-bridge.sh
Restart=on-failure
RestartSec=2s
